# debian/pictor.cron.daily
#!/bin/sh

run-one sudo -u www-data /usr/bin/pictor-thumbs
